# _system/deploy/com.myvilla.panel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.myvilla.panel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ivogiuliani/Code/myvilla-la/_system/scripts/approve.py</string>
        <string>--no-browser</string>
        <string>--port</string>
        <string>8787</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ivogiuliani/Code/myvilla-la</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/ivogiuliani/Code/myvilla-la/_system/logs/panel.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ivogiuliani/Code/myvilla-la/_system/logs/panel.log</string>
</dict>
</plist>
